# php/supervisor/conf.d/swoole.conf
[program:swoole]
command=php think swoole start
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/swoole.error.log
stdout_logfile=/var/log/supervisor/swoole.log